# envoy-exporter.service
[Unit]
Description=Envoy Prometheus Exporter
Documentation=https://github.com/John-MustangGT/envoy-exporter
After=network.target
Wants=network.target

[Service]
Type=simple
User=envoy-exporter
Group=envoy-exporter
WorkingDirectory=/var/lib/envoy-exporter

# Main command
ExecStart=/usr/local/bin/envoy-exporter -config /etc/envoy-exporter/envoy_config.xml

# Restart policy
Restart=always
RestartSec=10
TimeoutStopSec=30

# Logging configuration
StandardOutput=journal
StandardError=journal
SyslogIdentifier=envoy-exporter

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=yes
ReadWritePaths=/var/lib/envoy-exporter

# Environment variables (optional)
# Environment=DEBUG=false
# Environment=DISABLE_LOG_TIMESTAMP=true
# Environment=LOG_PREFIX=[solar-monitor]

[Install]
WantedBy=multi-user.target
